# linux-client/olv-vpn.service
[Unit]
Description=OpenLay VPN client (reconnect-on-TTL daemon)
Documentation=https://github.com/openlay/openlay-vpn-client-release
After=network-online.target systemd-resolved.service
Wants=network-online.target systemd-resolved.service

[Service]
Type=simple
# Run the self-contained binary dropped by install.sh at /usr/sbin/olv-vpn
# (PyInstaller onefile — no separate Python runtime or module path needed).
ExecStart=/usr/sbin/olv-vpn daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Sandboxing — daemon needs root (for wg-quick, /etc/wireguard, TPM device)
# but doesn't need access outside these paths.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/etc/wireguard /etc/olv /var/lib/olv /run
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
